# wemo_gpl_publication_SNS-9041/wemo_gpl_sns_package-9041/package/WeMo_SNS/files/usr/lib/opkg/info/WeMo.list
/etc/init.d/wemo
/etc/icon.jpg
/etc/ver.txt
/etc/sensor.jpg
/sbin/startWemo.sh
/sbin/web/setup.xml
/sbin/wemoApp
